# bot.conf
[program:ydnhbot]
command=python ydnhbot.py 
directory=/path/to/
user=root